# spot-price/crontab
1 * * * * python /spot_price.py >> /logfile
