# myupway/crontab
*/5 * * * * python /poll_metrics.py >> /logfile
